# flask.service
[Unit]
Description=Flask Web Server for Service Dashboard
After=network.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/home/jmagar/dashboard/venv/bin/python3 /home/jmagar/dashboard/flask.py
Restart=on-failure
RestartSec=30s
User=jmagar
Group=jmagar
WorkingDirectory=/home/jmagar/dashboard
EnvironmentFile=/home/jmagar/dashboard/.env

# Security
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true
ReadOnlyPaths=/
ReadWritePaths=/home/jmagar/dashboard/logs /home/jmagar/dashboard/templates

# Network
RuntimeDirectory=flask
RuntimeDirectoryMode=755

# Resource limits
CPUQuota=25%
MemoryLimit=256M
TasksMax=100

[Install]
WantedBy=multi-user.target
